# scripts/launchd/ai.mnelo.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.mnelo.mcp</string>
    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PY__</string>
        <string>__LIVE_ROOT__/mcp_server.py</string>
        <string>--transport</string>
        <string>sse</string>
        <string>--host</string>
        <string>127.0.0.1</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HERMES_HOME</key>
        <string>__HERMES_HOME__</string>
        <!-- Port 显式走 env (config chain: env > toml > default).
             Plist 不传 --port 让 argparse 走 config.server_port. -->
        <key>HERMES_MEMORY_SERVER_PORT</key>
        <string>8086</string>
        <!-- 不显式 PATH, 让 launchd 从 VIRTUAL_ENV prepend venv bin.
             显式 PATH 任何本地 user 可以丢恶意 python3 shim → RCE. -->
        <key>VIRTUAL_ENV</key>
        <string>__VENV_DIR__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!-- 防止 recall 死循环拖垮 stack -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>2048</integer>
        <key>NumberOfProcesses</key>
        <integer>512</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
        <key>NumberOfProcesses</key>
        <integer>1024</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__HERMES_HOME__/logs/mnelo.mcp.log</string>
    <key>StandardErrorPath</key>
    <string>__HERMES_HOME__/logs/mnelo.mcp.error.log</string>
    <key>WorkingDirectory</key>
    <string>__LIVE_ROOT__</string>
</dict>
</plist>