# scripts/launchd/ai.mnelo.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.mnelo.mcp</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/apple/hermes-agent/venv/bin/python3</string>
        <string>/Users/apple/.hermes/memory/mcp_server.py</string>
        <string>--transport</string>
        <string>sse</string>
        <string>--host</string>
        <string>127.0.0.1</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MNELO_HOME</key>
        <string>/Users/apple/.hermes</string>
        <!-- [7/19] Port 显式走 env (config chain: env > toml > default).
             Plist 不传 --port 让 argparse 走 config.server_port. -->
        <key>MNELO_MEMORY_SERVER_PORT</key>
        <string>8086</string>
        <!-- [7/19 P1-7 fix] 去掉显式 PATH, 让 launchd 自动从 VIRTUAL_ENV prepend venv bin.
             之前的 PATH 包含 /Users/apple/hermes-agent/venv/bin (world-readable),
             任何本地 user 可以丢恶意 python3 shim → RCE on next kickstart. -->
        <key>VIRTUAL_ENV</key>
        <string>/Users/apple/hermes-agent/venv</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!-- [7/19 P2-10] ResourceLimits: 防止 recall 死循环拖垮 stack -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>2048</integer>
        <key>NumberOfProcesses</key>
        <integer>512</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
        <key>NumberOfProcesses</key>
        <integer>1024</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/apple/.hermes/logs/mnelo.mcp.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/apple/.hermes/logs/mnelo.mcp.error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/apple/.hermes/memory</string>
</dict>
</plist>